# upstart/btsync.conf
# btsync daemon
author "miphreal"

start on runlevel [2345]
stop on runlevel [016]

setuid miph
setgid miph

exec btsync --config /home/miph/.btsync --nodaemon
